# supervisor/tictacwar.conf
[program:tictacwar]
command=/usr/local/bin/racket /opt/tictacwar/rack-serve.rkt
directory=/opt/tictacwar
startsecs=5
stderr_logfile=/var/log/tictacwar.log
